# test/docker/mongo-replset/init-replicaset.conf
[program:init-replicaset]
; this program is expected to exit (with status 0).
; it waits for all replica set members to be available then uses the
; mongo shell to initiate the replica set.
command = bash -c '/tmp/check_mongod.sh && /usr/bin/mongo --host localhost --port 30001 </tmp/init.js'
priority = 999
autorestart = false
startsecs = 0
startretries = 0
